# garuda-kde-settings/etc/skel/.config/autostart/org.kde.latte-dock.desktop
[Desktop Entry]
Exec=/usr/bin/latte-dock
Hidden=false
